# usr/lib/systemd/system/cloudflare-dynamic-dns-client.timer
[Unit]
Description=Periodic Coudflare Dynamic DNS Client update

[Timer]
OnBootSec=1m
OnUnitInactiveSec=1h

[Install]
WantedBy=multi-user.target
